# systemd/backup.service
[Unit]
Description=Backup Home

[Service]
ExecStart="$HOME/dotfiles/script/backup.sh