# systemd/dotfiles-sync.service
[Unit]
Description=Sync dotfiles with git

[Service]
ExecStart="$HOME/dotfiles/scripts/sync.sh"
StandardOutput=""
StandardError=""