# systemd/dotfiles-sync.service
# systemd/dotfiles-sync.service

[Unit]
Description=Sync dotfiles with git every 30 minutes
Wants=dotfiles-sync.timer

[Service]
Type=simple
ExecStart=%h/dotfiles/scripts/sync.sh
StandardOutput=journal
StandardError=journal